# ram.conf
[general]
appname=scwal

[program:soda]
command=make soda_server

[program:plugins]
command=make plugins_with_sleep

[program:ramonahttpfend]
command=<httpfend>
listen=tcp://0.0.0.0:5588

